# supervisord.conf
[supervisord]
nodaemon=true

[program:ml_heating_backend]
command=python3 -m src.main
directory=/app
autostart=true
autorestart=true
stderr_logfile=/data/logs/ml_heating_backend.err.log
stdout_logfile=/data/logs/ml_heating_backend.out.log

[program:dashboard]
command=streamlit run dashboard/app.py --server.port=3001 --server.address=0.0.0.0 --server.headless=true --server.enableCORS=false --server.enableXsrfProtection=false --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
stderr_logfile=/data/logs/dashboard.err.log
stdout_logfile=/data/logs/dashboard.out.log

[program:health]
command=python3 /app/dashboard/health.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/data/logs/health.err.log
stdout_logfile=/data/logs/health.out.log
